# examples/streamerbrainz.service
[Unit]
Description=StreamerBrainz daemon
Documentation=https://github.com/yourusername/streamerbrainz
After=network.target

[Service]
Type=simple

# Optional but a good idea: hides the runtime creds dir from other services via a mount namespace
PrivateMounts=yes

# Load encrypted credential (create with: systemd-creds encrypt --name=plex-token - - > ~/.config/streamerbrainz/plex-token.cred)
LoadCredentialEncrypted=plex-token:%h/.config/streamerbrainz/plex-token.cred

# Path to the binary
# %d expands to the credential directory for this service
ExecStart=%h/.local/bin/streamerbrainz \
    -ir-device /dev/input/event6 \
    -camilladsp-ws-url ws://127.0.0.1:1234 \
    -ipc-socket %t/streamerbrainz.sock \
    -camilladsp-min-db -65.0 \
    -camilladsp-max-db 0.0 \
    -camilladsp-update-hz 100 \
    -webhooks-port 3001 \
    -plex-server-url http://plex.home.arpa:32400 \
    -plex-token-file=%d/plex-token \
    -plex-machine-id YOUR_MACHINE_IDENTIFIER_HERE \
    -log-level info

# Restart policy
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=streamerbrainz

## Monitoring user space service logs
```bash
journalctl _UID=$(id -u) _SYSTEMD_USER_UNIT=streamerbrainz.service -f
```

[Install]
WantedBy=default.target
