# examples/streamerbrainz.service
[Unit]
Description=StreamerBrainz daemon
Documentation=https://github.com/nikoskalogridis/streamerbrainz
After=network.target

[Service]
Type=simple

# Run StreamerBrainz using the default config location (~/.config/streamerbrainz/config.yaml)
# All configuration is in the YAML file
ExecStart=%h/.local/bin/streamerbrainz

# Restart policy
Restart=on-failure
RestartSec=5s

# Resource limits
LimitNOFILE=65536

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=streamerbrainz

## Setup instructions:
# 1. Copy this file to ~/.config/systemd/user/streamerbrainz.service
# 2. Copy examples/config.yaml to ~/.config/streamerbrainz/config.yaml and edit it
# 3. Store your Plex token in a file:
#    mkdir -p ~/.config/streamerbrainz
#    echo -n "YOUR_PLEX_TOKEN" > ~/.config/streamerbrainz/plex-token
#    chmod 600 ~/.config/streamerbrainz/plex-token
# 4. In your config.yaml, set:
#    plex:
#      enabled: true
#      token_file: ~/.config/streamerbrainz/plex-token
#      server_url: http://plex.home.arpa:32400
#      machine_id: YOUR_MACHINE_IDENTIFIER
# 5. Reload and start:
#    systemctl --user daemon-reload
#    systemctl --user enable --now streamerbrainz

## Monitoring user space service logs
# journalctl _UID=$(id -u) _SYSTEMD_USER_UNIT=streamerbrainz.service -f

[Install]
WantedBy=default.target
